# src/resources/macos/com.aspauldingcode.wawona.compositor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aspauldingcode.wawona.compositor</string>
    
    <key>ProgramArguments</key>
    <array>
        <!-- UPDATE THIS PATH to your actual compositor binary -->
        <string>/path/to/build/Wawona</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/tmp/wawona.log</string>
    
    <key>StandardErrorPath</key>
    <string>/tmp/wawona.error.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>WAYLAND_DISPLAY</key>
        <string>wayland-0</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    
    <!-- Optional: ThrottleInterval prevents rapid restarts -->
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
